# payload/usr/local/mhealth/templates/com.mhealth.activity.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mhealth.activity</string>
    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON__</string>
        <string>/usr/local/bin/mhealth-activity</string>
    </array>
    <key>StartInterval</key>
    <integer>60</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/mhealth-activity.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/mhealth-activity.stderr.log</string>
</dict>
</plist>
